# layer/meta-application/recipes-ssl/event-sensors/files/event_sensors.timer
[Install]
WantedBy=basic.target

[Timer]
OnBootSec=2min
